# systemd/system/disk-scrub.timer
[Unit]
Description=timer for server update

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target