# systemd/system/server-motd-update.timer
[Unit]
Description=timer for server update

[Timer]
OnBootSec=1min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target