# systemd/system/server-motd-update.timer
[Unit]
Description=timer for server update

[Timer]
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target