# systemd/system/server-update.timer
[Unit]
Description=timer for server update

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target