# systemd/system/system-update.timer
[Unit]
Description=timer for system updates
[Timer]
OnUnitInactiveSec=24h
[Install]
WantedBy=timers.target