# systemd_old/system/server-update.timer
[Unit]
Description=Timer for Update Server Service

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target