# systemd_old/backup/server-backup.timer
[Unit]
Description=Timer for Backup Service

[Timer]
OnUnitInactiveSec=6h

[Install]
WantedBy=timers.target