# systemd/docker/docker-healthcheck.timer
[Unit]
Description=Timer for Docker Healthcheck Service

[Timer]
OnUnitInactiveSec=10m

[Install]
WantedBy=timers.target