# systemd/docker/docker-healthcheck.timer
[Unit]
Description=Docker healthcheck timer.

[Timer]
OnUnitInactiveSec=10m

[Install]
WantedBy=timers.target