# systemd/backup/backup-movies.timer
[Unit]
Description=Timer for syncing shows to backup host.

[Timer]
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target