# systemd/backup/media-backup.timer
[Unit]
Description=Timer for syncing files to backup host.

[Timer]
OnUnitInactiveSec=6h

[Install]
WantedBy=timers.target