# system/systemd/system/foss-motdupdate.timer
[Unit]
Description=Timer for Motd Update Service
[Timer]
OnBootSec=1min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target