# system/systemd/system/foss-docker-clean.timer
[Unit]
Description=Timer for Docker Cleanup Service

[Timer]
OnUnitInactiveSec=12h

[Install]
WantedBy=timers.target