# system/systemd/system/foss-diskscrub.timer
[Unit]
Description=Timmer for Disk Scrubber Service

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target