# system/systemd/docker/fods-media.target
[Unit]
Description=FlanneryOnline Media Services
ConditionPathExistsGlob=/opt/server/bin/docker/media/*
After=multi-user.target