# systemd/docker/containers/tautulli.service
[Unit]
Description=Service running tautulli in Docker.
Requires=docker-create.service plex.service
After=docker-create.service plex.service

[Service]
ExecStartPre=/usr/bin/docker start tautulli
ExecStart=/usr/bin/docker wait tautulli
ExecStop=/usr/bin/docker stop tautulli

[Install]
WantedBy=multi-user.target